# .config/systemd/user/dotfiles-sync.service
[Unit]
Description=Sync dots

[Service]
ExecStart=/usr/bin/zsh -c 'SERVICE=1 dots-sync'
Restart=on-failure

[Install]
WantedBy=default.target
